# deploy/systemd/estate-sale-scanner.timer
[Unit]
Description=Run Estate Sale Finder weekly

[Timer]
OnCalendar=Wed *-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
