# deploy/systemd/estate-sale-image-pull.timer
[Unit]
Description=Pull Estate Sale Finder scanner image hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
