# deploy/systemd/estate-sale-scanner.timer
[Unit]
Description=Run Estate Sale Finder daily

[Timer]
OnCalendar=*-*-* 06:15:00
Persistent=true

[Install]
WantedBy=timers.target
