# deploy/systemd/estate-sale-scanner.service
[Unit]
Description=Run Estate Sale Finder scanner
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/estate-sale-finder
ExecStartPre=/usr/bin/docker compose run --rm scanner migrate
ExecStart=/usr/bin/docker compose run --rm scanner run
